# srv/production/super_askreddit_but_ai.conf
[program:askreddit_but_ai]
; Command to start app
command = /usr/local/src/askreddit_but_ai/srv/production/gunicorn_start.bash
user = arai_user                                                          ; User to run as
stdout_logfile = /usr/local/src/askreddit_but_ai/logs/gunicorn_supervisor.log   ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                       ; Set UTF-8 as default encoding
